# core/services/mavlink_router/run
#!/usr/bin/env bash
DIR=$(dirname "$(readlink -f "$0")")
$DIR/install/usr/bin/mavlink-routerd "$@"